# config/dapper/autostart/blueman-applet.desktop
[Desktop Entry]
Type=Application
Version=1.0
Path=/home/coderobe
Name=blueman-applet
Exec=blueman-applet
