# config/dapper/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Version=1.0
Path=/home/coderobe
Name=nm-applet
Exec=nm-applet
