# config/dapper/autostart/plank.desktop
[Desktop Entry]
Type=Application
Version=1.0
Path=/home/coderobe
Name=plank
Exec=plank
