# config/dapper/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Version=1.0
Path=/home/coderobe
Name=redshift
Exec=redshift-gtk
